# packaging/systemd/hyperion-agent.service
[Unit]
Description=hyperion agent
After=network-online.target nginx.service mariadb.service postgresql.service
Wants=network-online.target

[Service]
Type=simple
ExecStart=/usr/sbin/hyperion-agent --config /etc/hyperion/agent.toml
Restart=on-failure
RestartSec=3s

# Hardening
# ProtectSystem=true (not "full") because the agent legitimately writes
# /etc/passwd, /etc/nginx/sites-available, /etc/php/*/fpm/pool.d,
# /etc/systemd/system/*, etc. as part of provisioning hostings.
# ProtectSystem=true keeps /usr and /boot read-only (which the agent
# never needs to write to) while letting useradd & friends do their job.
NoNewPrivileges=true
ProtectSystem=true
ProtectHome=false
ProtectKernelTunables=true
ProtectKernelModules=true
PrivateTmp=true
LogsDirectory=hyperion
RuntimeDirectory=hyperion
LimitNOFILE=65536

[Install]
WantedBy=multi-user.target
